# systemd/mtga-companion-staging.service
[Unit]
Description=MTGA Companion BFF Staging
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mtga-companion
Group=mtga-companion
ExecStart=/usr/local/bin/mtga-bff-staging
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mtga-bff-staging

# Environment — staging-specific values
Environment=MTGA_ENV=staging
Environment=BFF_PORT=8081

# Environment file — populated by staging-deploy.yml from /vaultmtg/staging/* SSM parameters
EnvironmentFile=/etc/mtga-companion-staging/env

# Resource limits — staging shares the EC2 instance with production;
# caps prevent a runaway staging process from starving the prod service.
MemoryMax=512M
CPUQuota=50%

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target
