# systemd/mtga-companion.service
[Unit]
Description=MTGA Companion BFF Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mtga-companion
Group=mtga-companion
ExecStart=/usr/local/bin/mtga-bff
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mtga-bff

# Environment file — populated by setup.sh from Secrets Manager / SSM
EnvironmentFile=/etc/mtga-companion/env

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target
